# jira-rofi.timer
[Unit]
Description=Timer for the rofi jira cache service
After=network.target

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=default.target
